# contrib/systemd/user/kde-fix-screen-recording.service
[Unit]
Description=KDE Fix Screen Recording

[Service]
Type=simple
ExecStart=/usr/bin/kde-fix-screen-recording

[Install]
WantedBy=default.target
